# my_rule.rules
KERNEL=="sdb", ACTION=="add",  RUN+="/usr/bin/device_added.sh"
